# spring_2026_ID/pi/shared/etc/systemd/system/dnsmasq.service.d/wait-for-network.conf
[Unit]
After=NetworkManager-wait-online.service
Wants=NetworkManager-wait-online.service
