# spring_2026_ID/pi/shared/etc/systemd/system/orc-sensors.timer
[Unit]
Description=Poll environmental sensors

[Timer]
OnBootSec=60
OnUnitActiveSec=60
AccuracySec=5

[Install]
WantedBy=timers.target
