# spring_2026_ID/pi/shared/etc/systemd/system/orc-maintenance-check.service
[Unit]
Description=Check GitHub for station maintenance mode
After=network-online.target
Wants=network-online.target
Before=orc-api.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/orc-maintenance-check
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
