# spring_2026_ID/pi/shared/etc/systemd/system/orc-sensors.service
[Unit]
Description=Log environmental sensor readings
After=local-fs.target

[Service]
Type=oneshot
User=pi
ExecStart=/usr/local/bin/orc-sensors
TimeoutStartSec=30
